# scripts/bithumb-bot.service
[Unit]
Description=Bithumb Auto Trading Bot v2
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
WatchdogSec=600
User=bythejune
Group=bythejune
WorkingDirectory=/home/bythejune/projects/bithumb-bot-v2
ExecStart=/home/bythejune/projects/bithumb-bot-v2/venv/bin/python run_bot.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# 로그 (journald 자동 로테이션)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bithumb-bot

# 리소스 제한
MemoryMax=2G
CPUQuota=80%

# 보안
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/bythejune/projects/bithumb-bot-v2/data

[Install]
WantedBy=multi-user.target
